# config/systemd/user/notmuch-new.timer
[Unit]
Description=Notmuch mail indexing and retrieval service timer

[Timer]
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
